# site-cookbooks/reuters-demo/templates/default/reuters-solr-proxy.conf.erb
# Upstart script for Nodejs proxy server
# Generated by Chef for <%= node[:fqdn] %>
description "Node.js security proxy for Solr"

start on filesystem and started networking
stop on shutdown

respawn

script
  exec 2>> /var/log/upstart/reuters-err.log
  exec 1>> /var/log/upstart/reuters-access.log
  solr-security-proxy --port <%= @proxy_port %> --backend.port <%= @solr_port %>
end script
