# site-cookbooks/reuters-demo/templates/default/reuters-solr-proxy.conf.erb
# Generated by Chef for <%= node[:fqdn] %>
description "Node.js security proxy server for AJAX-SOLR Reuters Demo"
 
start on filesystem and started networking
stop on shutdown
 
# forever forks the Node.js process
expect fork
 
env NODE_BIN_DIR="<%= node['reuters-demo']['node_bin'] %>"
# Set the NODE_PATH to the Node.js main node_modules directory.
env NODE_PATH="<%= node['reuters-demo']['node_path'] %>"

script
  PATH=$NODE_BIN_DIR:$PATH
  exec sudo -u <%= @user %> forever -a --sourceDir <%= node['reuters-demo']['proxy_root']%> start -l <%= @forever_logfile %> <%= @node_appfile %> >> <%= @forever_logfile %>
end script

pre-stop script
    PATH=$NODE_BIN_DIR:$PATH
    exec sudo -u <%= @user %> forever stop <%= @node_appfile %> >> <%= @forever_logfile %>
end script
